# lid-handler.service
[Unit]
Description=Lid Handler

[Service]
Type=simple
ExecStart=lid-handler

[Install]
WantedBy=default.target
